# systemd/macnix-vm.service
[Unit]
Description=MacNix macOS Virtual Machine
After=network-online.target libvirtd.service macnix-gpu-setup.service
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/opt/macnix/scripts/launch-vm.sh
ExecStop=/bin/kill -TERM $MAINPID
Restart=on-failure
RestartSec=10
KillMode=mixed
TimeoutStopSec=30
# CPU and memory limits
LimitMEMLOCK=infinity
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
